# s6/etc/s6-overlay/s6-rc.d/redis/run
#!/usr/bin/env sh
redis-server --loglevel nothing