# s6/etc/s6-overlay/s6-rc.d/redis/run
#!/bin/sh
redis-server --loglevel nothing