# s6/etc/s6-overlay/s6-rc.d/worker/run
#!/usr/bin/env sh
uv run --quiet --directory /app/ taskiq worker terebii.app:broker --workers 1 --no-configure-logging