# s6/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv sh
uv run --directory /app/ taskiq worker terebii.app:broker --workers 1 --no-configure-logging