# s6/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/with-contenv sh
uv run --quiet --directory /app/ taskiq scheduler terebii.app:scheduler --no-configure-logging